# lib/systemd/freshy.timer.template
[Unit]
Description=Freshy periodic sync timer

[Timer]
Persistent=true
OnBootSec=2min
OnUnitActiveSec={{.Schedule}}
AccuracySec=30s
Unit=freshy.service

[Install]
WantedBy=timers.target
