# lib/systemd/freshy.service
[Unit]
Description=Freshy package updater (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/bin/freshy sync
Environment=HOME=%h
Nice=10
